# etc/forge.service
[Unit]
Description=Forge

[Service]
Restart=always
ExecStart=forge -host ${HOSTNAME}

[Install]
WantedBy=multi-user.target
